# factory_installer/factory_ui.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# description "Factory installer service"
# author "chromium-os-dev@googlegroups.com"
#
# when dump-boot-stats has run, start up factory installer.

start on started factory_install

exec tail -F /var/log/softwareupdate.log > /dev/tty1
